# roles/node_exporter/templates/nodeexporter.service.j2
Description=Nodeexporter
After=network.target
 
[Service]
User={{ user }}
Group={{ group }}
Type=simple
ExecStart=/usr/local/bin/node_exporter
 
[Install]
WantedBy=multi-user.target